# tools/simple-http-server/app.service
[Unit]
Description=Dtlab login

[Service]
ExecStart=/root/app
